# src/openwrt/otbr
#!/bin/sh /etc/rc.common

START=90

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/otbr-agent -d7 -v /dev/ttyUSB0 115200
	procd_close_instance
}
