# third_party/wpantund/wpantund.service.in
[Unit]
Description=wpantund
After=network.target
ConditionPathExists=@sbindir@/wpantund

[Service]
EnvironmentFile=-@sysconfdir@/default/wpantund
ExecStart=@sbindir@/wpantund $WPANTUND_OPTS
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=SIGKILL

[Install]
WantedBy=multi-user.target
Alias=wpantund.service
