# src/agent/otbr-agent.service.in
[Unit]
Description=Border Router Agent
After=wpantund.service
ConditionPathExists=@sbindir@/otbr-agent

[Service]
EnvironmentFile=-@sysconfdir@/default/otbr-agent
ExecStart=@sbindir@/otbr-agent $OTBR_AGENT_OPTS
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=SIGKILL

[Install]
WantedBy=multi-user.target
Alias=otbr-agent.service
